# periodic/daily/902.my-timemachine-status
#!/bin/sh

if [ -r /etc/defaults/periodic.conf ]
then
	. /etc/defaults/periodic.conf
	source_periodic_confs
fi

TMSTATUS="$( /usr/libexec/PlistBuddy -c "print AutoBackup" /Library/Preferences/com.apple.TimeMachine.plist  2> /dev/null )"

if [ "${TMSTATUS:-true}" != "true" ]; then
	daily_my_timemachine_state=NO
fi

echo ""
case "${daily_my_timemachine_state:-NO}" in
    [Yy][Ee][Ss])
		echo "Time Machine Backup for this Mac Status: AutoBackup Enabled"
		echo "The last 10 snapshots:"
		echo "---------------------------------------------------"
		tmutil listbackups | tail | awk '{print "    "$0}'
		rc=$?
		;;
	*)
		echo "Time Machine Backup for this Mac Status: AutoBackup Disabled"
		rc=0
		;;
esac

exit $rc
