# periodic/weekly/901.certbot-renew
#!/bin/bash

if [ -r /etc/defaults/periodic.conf ]; then
	. /etc/defaults/periodic.conf
	source_periodic_confs
fi

if [ ! -x /usr/local/bin/certbot ]; then
	certbot_renew=NO
fi

case "${certbot_renew:-NO}" in
    [Yy][Ee][Ss])
        if [ -x  /usr/local/sbin/mailc ]; then
            /usr/local/bin/certbot renew --post-hook "/usr/local/sbin/mailc restart all" 2>&1
        else
            /usr/local/bin/certbot renew 2>&1
        fi
		rc=$?
		;;
	*)
		rc=0
		;;
esac

exit $rc
