# periodic/weekly/901.certbot-renew
#!/bin/bash

if [ -r /etc/defaults/periodic.conf ]; then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

if [ ! -x /usr/local/bin/certbot ]; then
    certbot_renew=NO
fi

case "${certbot_renew:-NO}" in
    [Yy][Ee][Ss])
        /usr/local/bin/certbot renew 2>&1
        rc=$?
        ;;
    *)
        rc=0
        ;;
esac

exit $rc
